# deployment/teams-notetaker-web.service
[Unit]
Description=Teams Meeting Transcript Summarizer - Web Dashboard
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=/home/%u/projects/teams-notetaker
Environment="PATH=/home/%u/projects/teams-notetaker/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Run web server
ExecStart=/home/%u/projects/teams-notetaker/venv/bin/python -m src.main serve --port 8001

# Restart on failure
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=teams-notetaker-web

# Resource limits
MemoryLimit=512M

[Install]
WantedBy=default.target
